# packages/agent/systemd/deployer-agent.service

[Unit]
Description="deployer"
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/opt/deployer/agent start
EnvironmentFile=/opt/deployer/.env
Restart=always
PrivateTmp=yes
User=root
