# packages/agent/systemd/tech-stack-dev-api-blue@.service

[Unit]
Description="tech-stack-dev-api-#%i"
After=network.target

[Install]
WantedBy=multi-user.target

[Service]
ExecStart=/opt/deployer/apps/tech-stack-dev/packages/api
EnvironmentFile=/opt/deployer/apps/tech-stack-dev/packages/api
Restart=always
DynamicUser=yes

RemainAfterExit=yes
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
RestrictNamespaces=yes
LockPersonality=yes
PrivateDevices=yes
ProtectClock=yes
ProtectKernelLogs=yes
RestrictSUIDSGID=yes
ProtectHome=false
PrivateUsers=yes
RuntimeDirectory="tech-stack-dev-api-#%i"
StateDirectory="tech-stack-dev-api-#%i"
CacheDirectory="tech-stack-dev-api-#%i"
LogsDirectory="tech-stack-dev-api-#%i"
ConfigurationDirectory="tech-stack-dev-api-#%i"
